# rtc/supervisor/conf/spider-8000.conf
[program:spider-8000]
command=/home/admin/venv/spider/bin/python /home/admin/spider/server.py --config=/home/admin/spider/spiderd.conf --port=8000
stderr_logfile = /var/log/supervisor/spider-stderr.log
stdout_logfile = /var/log/supervisor/spider-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10